# scripts/trade-analyzer-api.service
[Unit]
Description=Trade Strategy Analyzer API Server
After=network.target

[Service]
Type=simple
User=alvin
WorkingDirectory=/home/alvin/.openclaw/workspace/trade_strategy_analyzer/scripts
ExecStart=/usr/bin/python3 /home/alvin/.openclaw/workspace/trade_strategy_analyzer/scripts/api_server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
